# services/ansible-auto-update.timer
[Unit]
Description=Run Ansible Auto Update every day

[Timer]
OnCalendar=daily
Unit=ansible-auto-update.service

[Install]
WantedBy=timers.target
